# obsdfreqd.rc
#!/bin/ksh

daemon="/usr/local/sbin/obsdfreqd"

. /etc/rc.d/rc.subr

pexp="${daemon}.*"
rc_reload=NO
rc_bg=YES

rc_cmd $1
